# systemd/wifi-check.service
[Unit]
Description=WiFi Connection Monitor
Documentation=https://github.com/mokcontoro/headless_wifi_connect
After=network.target NetworkManager.service wifi-portal.service
Wants=NetworkManager.service wifi-portal.service

[Service]
Type=simple
User=root
ExecStart=/bin/bash /usr/local/bin/check_wifi.sh
Restart=on-failure
RestartSec=30
StartLimitBurst=5
StartLimitIntervalSec=300
StandardOutput=journal
StandardError=journal

# Security
NoNewPrivileges=false

# Prevent service from interfering during network changes
TimeoutStopSec=10

[Install]
WantedBy=multi-user.target
